# deploy/systemd/photobooth.service
# /etc/systemd/system/photobooth.service
#
# Service kiosk de la borne photo. Démarre seul, plein écran, sur tty1,
# et se relance automatiquement en cas de crash (principal levier de stabilité).
#
# Activation (une fois, par le mainteneur) :
#     sudo cp deploy/systemd/photobooth.service /etc/systemd/system/
#     sudo systemctl daemon-reload
#     sudo systemctl enable photobooth.service
#     # libérer tty1 (sinon un curseur de login apparaît derrière le kiosk) :
#     sudo systemctl set-default multi-user.target
#     sudo systemctl disable getty@tty1.service
#     sudo systemctl mask getty@tty1.service
#
# Logs : journalctl -u photobooth -f

[Unit]
Description=Photobooth kiosk
After=multi-user.target photobooth-provision.service
Conflicts=getty@tty1.service
ConditionPathExists=!/boot/firmware/photobooth/disable-kiosk

[Service]
Type=simple
User=pi
SupplementaryGroups=gpio i2c video input render
WorkingDirectory=/home/pi/photobooth
ExecStart=/home/pi/photobooth/Photobooth.App --drm

# --- Rendu GPU DRM/KMS (vc4-kms-v3d requis dans /boot/firmware/config.txt) ---
# StartLinuxDrm passe par le GPU VideoCore IV via Mesa vc4 : transforms et opacité
# sont accélérés matériellement, ce qui rend les animations de cartes fluides.
# Validé sur Raspberry Pi OS Trixie — si la carte affiche un écran noir, repasser
# à --fbdev et remettre LIBGL_ALWAYS_SOFTWARE=1 + GALLIUM_DRIVER=llvmpipe.

# Plafonne le tas .NET (~384 Mo) pour ne pas faire swapper le Pi (1 Go).
Environment=DOTNET_GCHeapHardLimit=0x18000000

Restart=always
RestartSec=3
StandardInput=tty
TTYPath=/dev/tty1
TTYReset=yes
TTYVTDisallocate=yes

[Install]
WantedBy=multi-user.target
